# ja/Ports/dreamcast/howto.list
<!DOCTYPE HTML PUBLIC "-//W3C//DTD HTML 4.0 Transitional//EN">

<html>
<head>
<!-- Copyright (c) 1994-2003
        The NetBSD Foundation, Inc.  ALL RIGHTS RESERVED. -->
<meta http-equiv="Content-Type" content="text/html; charset=ISO-2022-JP">
<link rev="made" href="mailto:www@JP.NetBSD.ORG">
<title>NetBSD/dreamcast: How to use NetBSD/dreamcast</title>
</head>
<body bgcolor="#ffffff" text="#000000" link="#2020cc" vlink="#700080">

<h1>NetBSD/dreamcast の使い方</h1>

<i>&quot;NetBSD を動かすために Dreamcast を買ったところです。どうしたらいい?&quot</i>
<p>

何が問題なのかさえわかっていれば、 Sega Dreamcast で NetBSD を動くように
するのはいたって簡単なことです。このドキュメントは、 DC/NetBSD 初心者に加えて
NetBSD 初心者も対象としています。ただし、ある程度の UNIX 的な知識は
前提としています。このドキュメントに載せた情報は、
<a href="http://mail-index.netbsd.org/port-dreamcast/">port-dreamcast
メーリングリストのアーカイブ</a>、 Marcus Comstedt のホームページ <a
href="http://mc.pp.se/dc/">http://mc.pp.se/dc/</a> 、 BeTips.net <a
href="http://www.betips.net">http://www.betips.net</a> 、
<a href="http://members.wri.com/josht/dreamcast/dreamcast.html">Josh Tolbert's 
Dreamcast How-to</a>、および筆者個人の
試行錯誤の結果から選りすぐりました。
<p>

システムをまともに使うためには書き込み可能なファイルシステムが必要なので、
この HOWTO では Dreamcast をネットブートさせることに集中します。 CD
からのブートはテスト用には便利ですし、
「ただのお遊び」でそうしたい人たちもいるので、
CD からのブートに関する簡単な説明もしておきます。
<p>

<h3>1) 必要なもの</h3>

下記のものが必要です:
<p>

<ul>
  <li>Dreamcast ゲーム機 (2000 年 9 月以前に製造されたものをおすすめします
    (それより後に製造されたものでも動いたとの情報もありますので、
    メーリングリストのアーカイブで確認してください))
    <p>
  <li>Dreamcast キーボード
    <p>
  <li>オプション - Dreamcast マウス (本稿執筆時点ではドライバーがありません。
    あなたが書いてみませんか?)
    <p>
  <li>Windows が動いている PC
    <p>
  <li>NetBSD が動いている PC (Linux や他の *BSD でもよいはずです。
    BeOS は step 3 が行なえません)
    <p>
  <li>ブランク CD-R (CD-RW は Dreamcast では読めません)
    <p>
  <li><a  href="http://members.wri.com/josht/dreamcast/dreamcast-kernel.tar.gz">
    Josh のサイト</a> (または <a 
    href="http://www.schnarff.com/netbsd/dreamcast-kernel.tar.gz">ミラー</a>)
    で配布されているブート可能な NetBSD カーネル。
    <p>
  <li>"mkhybrid" と CD を焼くためのプログラム (筆者は BeOS で CD Burner
    を使っていますが、ふつうは "cdrecord" でしょう)
</ul>

<h3>2) ブータブルカーネル CD を作る</h3>

<ul>
  <li>ダウンロードしたカーネルを適当な場所に展開します。 audio.raw と
      kernel.raw というファイルがあるはずです。
    <p>
  <li>最初に audio.raw をオーディオトラックとして焼きます。残りのセッション用に
      ディスクをオープンしたままにしておいてください。 cdrecord
      では以下のようにします:
    <p><pre>	# cdrecord dev=0,1,0 -multi -audio -audio.raw</pre>
    <p>もちろん、この dev=0,1,0 は、お使いのシステムに合った値に変える必要が
       あります。これは cdrecord -scanbus で調べられます。
    <p>
  <li>kernel.raw をデータとして焼きます。 cdrecord では以下のようにします:
    <p><pre>	# cdrecord dev=0,1,0 -multi -xa1 netbsd-kernel.raw</pre>
</ul>
<p>
テスト CD を焼く場合:
<ul>
  <li>ダウンロードしたイメージを DiscJuggler を使って焼きます。この CD を
      Dreamcast に入れて電源を入れた際にブートメッセージが出れば、その
      Dreamcast は使えるということです; もしそのマシンが 2002 年 9 月より後に
      製造されたものである場合は、 <a
      href="http://www.netbsd.org/cgi-bin/subscribe_list.pl?list=port-dreamcast">
      メーリングリスト</a>でお知らせください。
  <p>
  <li>CD 上にファイルシステムを作ってテストしたい場合:
    <ul>
      <li><a href="http://www.netbsd.org">http://www.netbsd.org</a> から
          HPCSH tarball を取ってきて、適当な場所に置きます
      <p>
      <li>/home/dreamcast/dev/ で<pre>	# sh MAKEDEV all</pre> を実行します
      <p>
      <li>この配布物を CD に焼くため、 mkhybrid でイメージを作ります:
      <pre>	# mkhybrid -o DCBSD.img -V DCNetBSD -a -r /home/dreamcast/</pre>
      -V はボリュームラベル、 -a はこのディレクトリー以下の全ファイル、 -r は
      Rock Ridge 拡張 (必須ではないが、指定しても問題はない) を意味します。
      <p>
      <li>DiscJuggler CD を Dreamcast に入れて電源を入れます。ルート FS
          を尋ねるプロンプトでは "gdrom" と打ち込みます; ダンプデバイスを
          尋ねられたら "Enter" を押します。
      <p>
      <li>ブート CD を、さきほど作ったファイルシステム CD に入れ換えます
      <p>
      <li>ファイルシステムを尋ねるプロンプトで "cd9660" と打ち込みます;
          これでシングルユーザーモードに入るはずです。
    </ul>
</ul>

<h3>3) NFS サーバー & ファイルシステムの準備</h3>

<ul>
  <li>Josh のサイトから<a 
      href="http://members.wri.com/josht/dreamcast/dreamcast-world.tar.gz">
      world ファイル</a>を取ってきます。
    <p>
  <li>このファイルを、お使いの *BSD/Linux PC の適当な場所に展開します。
      
    <p>
  <li>/etc/exports ファイルを用意します (書式については exports(5) マニュアルを
      参照) 。 Josh のサイトに <a href="http://members.wri.com/josht/dreamcast/exports.txt">
      例</a>があります。以下の点に注意してください:
    <ul>
      <li>"-maproot=0" します。さもないと、 Dreamcast はファイルシステムを
          マウントできません。
      <li>共有可能なネットワークかホストを指定します; さもないと、
          全世界から書き込み可能になってしまい、
          危険です。
    </ul>
    <p>
  <li>mountd と nfsd が動いていることを確認します; 動いていれば、 "ps aux | grep 
      '(mount|nfs)'d" でそれらが表示されます。表示されない場合、
      それらを起動するため、 root で "mountd" と "nfsd" を実行します;
      今後のブート時にこれらを自動的に起動するようにするため、 /etc/rc.conf
      を編集して、 nfs_server=YES、 lockd=YES、 statd=YES を書きます。
    <p>
  <li>/etc/dhcdp.conf を設定します (このファイルの書式は man dhcpd.conf(5)
      参照) 。 <a href="http://members.wri.com/josht/dreamcast/dhcpd.conf">
      設定例</a>が Josh のサイトにあります。このほか <a 
      href="http://www.schnarff.com/netbsd/dhcpd.conf"> Alex Kirk のサイト</a>
      にも例があります。 dhcpd.conf を書くときには以下の点に注意してください:
    <ul>
      <li>host セクションに Dreamcast の MAC アドレスを指定 *しなければ
          なりません* 。 MAC アドレスは BBA を見ればわかります;
          ステッカーに書いてあります。
      <li>root-path オプションは、 
          world ファイルを展開した場所にする必要があります。
      <li>domain-name-servers などのオプションは厳密には必要ありませんが、
          書いておいたほうが簡単になります。
    </ul>
    <p>
  <li>dhcpd が動いていることを確認します; 動いていなければ <pre>#	sh /etc/rc.d/dhcpd start</pre>
      します。いつでも DHCP が動くようにするため、 /etc/rc.conf を編集して
      dhcpd=YES を書いておきます。
    <p>
  <li>/home/dreamcast/dev/ で MAKEDEV スクリプトを実行します。これで、
      他の必要なデバイスとともにコンソールが作られます。
    <p>
    <pre>       # sh MAKEDEV all</pre>
    <p>
  <li>Dreamcast 用のスワップファイルを作ります。任意のサイズにできますが、
      経験上は 32-64MB が最適です (小さすぎると使い物になりませんし、
      大きすぎると保守のために過大な物理 RAM を占有します) 。
      スワップファイルは Dreamcast のルートとしてエクスポートする
      ディレクトリーの下に置くようにします。
    <p>
    <pre>	# dd if=/dev/zero of=swap bs=32k count=1024 </pre> (サイズが 32MB の場合)
    <p>
  <li>/home/dreamcast/etc/fstab (または Dreamcast のルートとなる場所の fstab
      ファイル) を編集し、適切なファイルシステムをマウントするようにします:
    <pre>	192.168.2.69:/home/alex/dreamcast/dc / nfs rw,auto 0 0
	/.swap none swap sw 0 0
	/kern /kern kernfs rw
	/proc /proc procfs rw 0 0
	/dev/gdrom0c /mnt/gdrom cd9660 ro,noauto 0 0
    </pre> 
    <p>
</ul>

<h3>4) やってみよう</h3>

<ul>
  <li>CD を Dreamcast に入れて電源を入れます。 CD が正しく作れていれば、
      見慣れた NetBSD のブートアップシーケンスが始まります。
      もし、 2000 年 9 月より後に製造された Dreamcast
      でこれに成功したら、 <a 
      href="http://www.netbsd.org/cgi-bin/subscribe_list.pl?list=port-dreamcast">
      メーリングリスト</a>でお知らせください。
    <p>
  <li>問題が起きた場合、 /etc/exports か /etc/dhcpd.conf がおかしいのかも
      しれません。まずはこれらにくだらない間違いがないか、再確認するとよいでしょう
      ─私たちもそういう間違いをするので、気にしないように。
      説明なしにエラーが出た場合は、
      /usr/src/sys/sys/errno.h を見ると簡単な説明があります。
    <p>
   <li>行き詰まったら、<a 
       href="http://mail-index.netbsd.org/port-dreamcast">メーリングリストの
       アーカイブを見ましょう</a>。これは役立ちます。検索したい場合は
       Google の Advanced Search を使うのが最善です; ドメインを
       mail-index.netbsd.org に絞り、検索語に "dreamcast"
       を含めるようにします。
    <p>
   <li>それでもうまくいかなければ、<a 
       href="http://www.netbsd.org/cgi-bin/subscribe_list.pl?list=port-dreamcast">
       メーリングリスト</a>に投稿してください。問題の詳細を書き、関連する
       エラーメッセージ/設定ファイル/ログファイルを (サイズが小さければ) 引用してください。
    <p>
</ul>

<h3>5) 完了!</h3>

カスタムカーネルをいろいろ作って遊びたい場合は、 Josh の <a
href="http://members.wri.com/josht/dreamcast/ipslave.html">IP Slave ページ</a>
にある説明を読むことをおすすめします。
このページでは、カーネルを作るたびに CD
を焼かなくてもすむ方法について、説明しています。

<hr>
<a href="./">Up to <em>NetBSD/dreamcast Port Page</em></a>
<hr>

<DOCLINK>

<hr>
<address>
  <small>
  (連絡先 - <a href="http://www.netbsd.org/cgi-bin/feedback.cgi">英語</a>,
       <a href="mailto:www@JP.NetBSD.ORG">日本語:
       www@JP.NetBSD.ORG</a>)<br>
  $NetBSD: howto.list,v 1.9 2003/05/15 23:02:53 hubertf Exp $<br>
  <a href="../../Misc/disclaimer.html">Copyright &copy; 1994-2003
  The NetBSD Foundation, Inc.  ALL RIGHTS RESERVED.</a>
  </small>
</address>

</body>
</html>
