# ja/Ports/dreamcast/howto.list
<!DOCTYPE HTML PUBLIC "-//W3C//DTD HTML 4.0 Transitional//EN">

<html>
<head>
<!-- Copyright (c) 1994-2003
        The NetBSD Foundation, Inc.  ALL RIGHTS RESERVED. -->
<meta http-equiv="Content-Type" content="text/html; charset=utf-8">
<link rev="made" href="mailto:www@NetBSD.ORG">
<title>NetBSD/dreamcast: How to use NetBSD/dreamcast</title>
</head>
<body bgcolor="#ffffff" text="#000000" link="#2020cc" vlink="#700080">

<h1>How to use NetBSD/dreamcast</h1>

<i>&quot;I Just Bought a Dreamcast to Run NetBSD.  What do I do now?&quot</i>
<p>

Getting NetBSD up and running on the Sega Dreamcast is pretty easy, as long as
you know what's involved.  This document is aimed at the NetBSD newbie as
well as the DC/NetBSD newbie.  However, a familiarity with some flavor of
UNIX is recommended.  Information for this document has been culled from
the <a href="http://mail-index.netbsd.org/port-dreamcast/">port-dreamcast
mailing list archive</a>, Marcus Comstedt's homepage at <a
href="http://mc.pp.se/dc/">http://mc.pp.se/dc/</a>, BeTips.net <a
href="http://www.betips.net">http://www.betips.net</a>, 
<a href="http://members.wri.com/josht/dreamcast/dreamcast.html">Josh Tolbert's 
Dreamcast How-to</a>, and personal trial and error.
<p>

This HOWTO will focus on netbooting your Dreamcast, since you need a writable
filesystem to do anything useful with the system. As booting from a CD can
be useful for test purposes -- and since some people might want to do it
"just for the fun of it" -- brief instructions on that setup are included as
well.
<p>

<h3>1) Requirements</h3>

You will need the following:
<p>

<ul>
  <li>Dreamcast console (it is recommended that you use a unit manufactured
    September 2000 or earlier (some have reported that newer ones work as
    well, so check the mailing list archive to be sure))
    <p>
  <li>Dreamcast keyboard
    <p>
  <li>OPTIONAL - Dreamcast mouse (as of this writing there is no driver for
    it.  Maybe you could write one?)
    <p>
  <li>a PC running Windows
    <p>
  <li>a PC running NetBSD (Linux or any other *BSD should work, too.
    BeOS does not work for step 3)
    <p>
  <li>Blank CD-R's (CD-RW's cannot be read by the Dreamcast)
    <p>
  <li>Bootable NetBSD kernel from <a 
    href="http://members.wri.com/josht/dreamcast/dreamcast-kernel.tar.gz">
    Josh's site</a> (or a <a 
    href="http://www.schnarff.com/netbsd/dreamcast-kernel.tar.gz">mirror</a>).
    <p>
  <li>"mkhybrid" and a CD burning program (I used CD Burner under BeOS, but
    you'll most likely use "cdrecord")
</ul>

<h3>2) Make a Bootable Kernel CD</h3>

<ul>
  <li>Extract the kernel to a convenient location. You should now have to
      files, audio.raw and kernel.raw.
    <p>
  <li>Burn the audio.raw track as audio first, making sure that you 
      leave the disk open for more sessions. Under cdrecord, this is done 
      with:
    <p><pre>	# cdrecord dev=0,1,0 -multi -audio -audio.raw</pre>
    <p>Obviously, you need to replace dev=0,1,0 with the appropriate
       device for your system, found with cdrecord -scanbus.
    <p>
  <li>Burn the kernel.raw file as data. This is done under cdrecord as:
    <p><pre>	# cdrecord dev=0,1,0 -multi -xa1 netbsd-kernel.raw</pre>
</ul>
<p>
To burn the test CDs:
<ul>
  <li>Burn the image you downloaded using DiscJuggler. If you get boot
      messages when you put this in the Dreamcast and power up, you know
      your system is ready to roll; please report this to <a
      href="http://www.netbsd.org/cgi-bin/subscribe_list.pl?list=port-dreamcast">
      the mailing list</a> if your machine was made after September 2002.
  <p>
  <li>If you wish to create a filesystem on CD and test it out:
    <ul>
      <li>Grab the HPCSH tarballs from <a href="http://www.netbsd.org">
          http://www.netbsd.org</a> and untar them in a convenient location
      <p>
      <li>Run <pre>	# sh MAKEDEV all</pre> in /home/dreamcast/dev/
      <p>
      <li>Use mkhybrid to burn a CD of this distribution:
      <pre>	# mkhybrid -o DCBSD.img -V DCNetBSD -a -r /home/dreamcast/</pre>
      -V is for volume label, -a means all files in the directory, and -r is for
      Rock Ridge extensions (not necessarily required, but not harmful).
      <p>
      <li>Put your DiscJuggler CD in the Dreamcast and power up. When prompted
          for the root FS, type "gdrom"; hit "Enter" when it asks
          about the dump device.
      <p>
      <li>Swap the boot CD for your freshly made filesystem CD
      <p>
      <li>Type in "cd9660" at the filesystem prompt; this should put
          you into single user mode.
    </ul>
</ul>
<p>

<h3>3) Prepare the NFS Server & Filesystem</h3>

<ul>
  <li>Fetch the <a 
      href="http://members.wri.com/josht/dreamcast/dreamcast-world.tar.gz">
      world files</a> from Josh's site.
    <p>
  <li>Untar them into a conveniently exportable directory on your *BSD/Linux
      PC. 
    <p>
  <li>Set up your /etc/exports (man exports(5) for more info on this file's
      syntax). An <a href="http://members.wri.com/josht/dreamcast/exports.txt">
      example</a> can be found on Josh's site. Make sure that you've:
    <ul>
      <li>Set "-maproot=0", otherwise the Dreamcast won't have permissions
          to mount the filesystem. 
      <li>Specified a network or host the share is accessible to; otherwise
          you could end up making it writable to the world, a Really
          Bad Thing.
    </ul>
    <p>
  <li>Make sure mountd and nfsd are running; a simple "ps aux | grep 
      '(mount|nfs)'d" will show if they're currently up. If they're not,
      run "mountd" and "nfsd" as root to start them now; to ensure they
      start at boot-time in the future, edit /etc/rc.conf so that 
      nfs_server=YES, lockd=YES, and statd=YES.
    <p>
  <li>Set up your /etc/dhcdp.conf (man dhcpd.conf(5) for more info on this file's
      syntax). An <a href="http://members.wri.com/josht/dreamcast/dhcpd.conf">
      example</a> can be found a Josh's site, or alternatlely at <a 
      href="http://www.schnarff.com/netbsd/dhcpd.conf"> Alex Kirk's site</a>.
      Things to keep in mind when configuring dhcpd.conf:
    <ul>
      <li>You MUST specify the Dreamcast's MAC Address in the host
          section. This can be found by looking at your BBA; it's on 
          the sticker.
      <li>Your root-path option should be wherever your world files were
          extracted.
      <li>Options such as domain-name-servers aren't strictly necessary,
          but tend to make life easier.
    </ul>
    <p>
  <li>Make sure dhcpd is running; turn it on with <pre>#	sh /etc/rc.d/dhcpd start</pre>
      if it's not on. To ensure that DHCP is running whenever your box is, edit 
      /etc/rc.conf so that dhcpd=YES.
    <p>
  <li>Run the MAKEDEV script in /home/dreamcast/dev/.  This will create 
      your console, as well as other needed devicesi.
    <p>
    <pre>	# sh MAKEDEV all</pre>
    <p>
  <li>Create a swapfile for the Dreamcast. You may pick any size you wish,
      but practical considerations say 32-64MB works best (much smaller and
      it's not helpful, much larger and it takes excessive physical RAM to
      maintain). Make sure it put it in the directory you're exporting as
      the Dreamcast's root. 
    <p>
    <pre>	# dd if=/dev/zero of=swap bs=32k count=1024 </pre> (for 32MB)
    <p>
  <li>Edit /home/dreamcast/etc/fstab (or wherever you've placed the Dreamcast
      root) to mount the appropriate filesystems:
    <pre>	192.168.2.69:/home/alex/dreamcast/dc / nfs rw,auto 0 0
	/.swap none swap sw 0 0
	/kern /kern kernfs rw
	/proc /proc procfs rw 0 0
	/dev/gdrom0c /mnt/gdrom cd9660 ro,noauto 0 0
    </pre> 
    <p>
</ul>

<h3>4) Try it Out</h3>

<ul>
  <li>Put your CD into the Dreamcast and power up. If you made the CD
      correctly, a familiar NetBSD bootup sequence will begin. If you've 
      done this with a Dreamcast manufactured after September 2000,
      please notify the <a 
      href="http://www.netbsd.org/cgi-bin/subscribe_list.pl?list=port-dreamcast">
      mailing list</a>.
    <p>
  <li>If you run into problems, chances are your /etc/exports or /etc/dhcpd.conf
      are messed up. You'll want to first double-check them for stupid mistakes
      -- we all make them, don't feel bad. If you've got a specific error
      message with no description, look in /usr/src/sys/sys/errno.h for a brief
      explanation.
    <p>
   <li>If you're really stuck, <a 
       href="http://mail-index.netbsd.org/port-dreamcast">browse the mailing list
       archives</a>. They're quite helpful. If you want to search them, your
       best bet is to use Google's Advanced Search; restrict the domain to 
       mail-index.netbsd.org and make sure that you include "dreamcast" in your
       query.
    <p>
   <li>That failing, post to the <a 
       href="http://www.netbsd.org/cgi-bin/subscribe_list.pl?list=port-dreamcast">
       list itself</a>. Please be detailed, and include relevant error messages/
       configuration files/logfiles (if they're small).
    <p>
</ul>

<h3>5) You're done!</h3>

If you're going to be playing around with custom-compiled kernels a lot, I'd
recommend following the instructions on Josh's <a 
href="http://members.wri.com/josht/dreamcast/ipslave.html">IP Slave page</a>,
which gives instructions on a setup where you won't have to burn a CD for
each fresh new kernel.

<hr>
<a href="./">Up to <em>NetBSD/dreamcast Port Page</em></a>
<hr>

<DOCLINK>

<hr>
<address>
  <small>
  <a href="http://www.netbsd.org/cgi-bin/feedback.cgi">(Contact us)</a>
  $NetBSD: howto.list,v 1.9 2003/05/15 23:02:53 hubertf Exp $<br>
  <a href="../../Misc/disclaimer.html">Copyright &copy; 1994-2003
  The NetBSD Foundation, Inc.  ALL RIGHTS RESERVED.</a>
  </small>
</address>

</body>
</html>
